# rootfs/config/supervisor/socat.conf
[program:socat]
command=socat TCP-LISTEN:9222,fork,reuseaddr TCP:127.0.0.1:19222
priority=10
autorestart=true
stderr_logfile=/var/log/socat.err.log